# hornet-supervisord.conf
; Supervisor configuration for Hornet

; Configuration for supervisord itself.
[supervisord]
logfile=/var/log/supervisord/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=error
pidfile=/var/run/supervisord.pid
;nodaemon=false
minfds=1024
minprocs=200
user=root
childlogdir=/var/log/supervisord/
; Configuration for hornet.
[program:hornet]
command=hornet -v
directory=/opt/vfs/
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/hornet.log
